# scripts/konnexcore.service
[Unit]
Description=KonnexCore FastAPI backend
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=@USER@
Group=@USER@
WorkingDirectory=@PROJECT_DIR@
Environment="PATH=@PROJECT_DIR@/venv/bin:/usr/bin:/bin"
EnvironmentFile=-@PROJECT_DIR@/.env
ExecStart=@PROJECT_DIR@/venv/bin/uvicorn api.main:app --host 127.0.0.1 --port 8000 --workers 2
Restart=always
RestartSec=3

# Sandboxing — the backend doesn't need world write access.
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=@PROJECT_DIR@
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=yes
LockPersonality=yes

# Conservative resource limits.
LimitNOFILE=4096

[Install]
WantedBy=multi-user.target
